# apticron.service
[Unit]
Description=Launch apticron to notify of packages pending an update

[Service]
Type=oneshot
ExecStart=/usr/sbin/apticron
