# debian/cron.daily
#!/bin/sh 

RANGE=3600
number=`od -vAn -N2 -tu4 < /dev/urandom`

number=`expr $number "%" $RANGE`
sleep $number
test -x /usr/sbin/apticron || exit 0
/usr/sbin/apticron
