# debian/cron.daily
#!/bin/sh 

RANGE=3600
number=$RANDOM

let "number %= $RANGE"
sleep $number
test -x /usr/sbin/apticron || exit 0
/usr/sbin/apticron
