# debian/cron.daily
#!/bin/sh

test -x /usr/sbin/apticron && /usr/sbin/apticron --cron
